# configs/local-configs/autostart/pcloud.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=pcloud
Comment=pcloudstartup script
Exec=/home/deck/opt/pcloud
StartupNotify=false
Terminal=false