# os-init/configs/etc_systemd_system/legionhub-check.service
[Unit]
Description=Check USB

[Service]
ExecStart=/home/deck/usr/os-init/configs/is-legion-plugged

[Install]
WantedBy=multi-user.target
