# configs/system/etc_udev_rules.d/90-egpu-compute-mode.rules
ACTION=="add", SUBSYSTEM=="thunderbolt", ATTR{authorized}=="0", RUN+="/usr/local/bin/rtx-gpu-compute-only.sh", ATTR{authorized}="1"
