# configs/system/etc_udev_rules.d/95-monitor-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/sbin/xrand-scale-deck"
