# configs/system/etc_systemd_system/legionhub-check.service
[Unit]
Description=Check USB
Before=sddm.service

[Service]
ExecStart=/usr/sbin/is-legion-plugged

[Install]
WantedBy=multi-user.target
