# configs/system/etc_systemd_system/reset-plasma.service
Description=Reset KDE Plasma

[Service]
ExecStart=/usr/bin/bash -c 'plasmashell --replace & disown'
