# configs/dot-config/autostart/arandr.sh.desktop
[Desktop Entry]
Exec=/home/deck/usr/configs/unlinked/arandr.sh
Icon=dialog-scripts
Name=arandr.sh
Path=
Type=Application
X-KDE-AutostartScript=true
