# configs/dot-config/autostart/kwalletd5.desktop
[Desktop Entry]
Exec=kwalletd5
Icon=
Name=kwalletd5
Path=
Terminal=False
Type=Application
