# configs/system/updater.service
[Unit]
Description=Daily system update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/updater.sh
