# configs/system/updater.timer
[Unit]
Description=Run system update daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
