# os-init/configs/etc_udev_rules.d/95-monitor-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/deck/.local/bin/xrand-scale-deck"
